# crontab
0 2 * * * /home/ubuntu/limpiar_logs.sh >> /home/ubuntu/cron_logs.txt 2>&1
